# scripts/claude-to-im.service
[Unit]
Description=Claude-to-IM Discord Bridge Daemon
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/ser8/.codex/skills/claude-to-im
ExecStart=/usr/bin/node dist/daemon.mjs
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Environment
EnvironmentFile=-/home/ser8/.codex/skills/claude-to-im/config.env

[Install]
WantedBy=default.target
