# scripts/claude-to-im.service
[Unit]
Description=Claude-to-IM Discord Bridge Daemon
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/.codex/skills/claude-to-im
ExecStart=/usr/bin/node dist/daemon.mjs
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Environment
EnvironmentFile=-%h/.claude-to-im/config.env
Environment=PATH=%h/.local/bin:%h/.npm-global/bin:/usr/local/bin:/usr/bin:/bin
Environment=HOME=%h

[Install]
WantedBy=default.target
